# root/etc/services.d/readarr/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc /app/readarr/bin/Readarr \
    -nobrowser -data=/config
