# root/etc/s6-overlay/s6-rc.d/init-readarr-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /run/readarr-temp

# permissions
lsiown -R abc:abc \
    /config \
    /run/readarr-temp
